# debian/postinst
ln -fs  /opt/tracker/tracker.sh /usr/bin/tracker
mkdir -p /opt/tracker/logs
systemctl daemon-reload
systemctl enable tracker
systemctl start tracker